# systemd/lotl-detector.service
[Unit]
Description=LOTL Detector - Living off the Land Detection and Prevention
Documentation=https://github.com/your-org/lotl-detector
After=network.target syslog.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root

# Working directory
WorkingDirectory=/opt/lotl-detector

# Main executable
ExecStart=/opt/lotl-detector/.venv/bin/python -m lotl_detector --config /etc/lotl/detector.yaml

# Restart policy
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

# Graceful shutdown
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

# Security hardening (while maintaining required capabilities)
NoNewPrivileges=no
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_NET_ADMIN CAP_BPF CAP_PERFMON CAP_SYS_PTRACE
AmbientCapabilities=CAP_SYS_ADMIN CAP_NET_ADMIN CAP_BPF CAP_PERFMON

# Resource limits
LimitNOFILE=65536
LimitMEMLOCK=infinity

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lotl-detector

# Environment
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target

